# packages/aur/linux-arctis-manager/arctis-manager.service
[Unit]
Description=Arctis Manager
StartLimitIntervalSec=1min
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/lam-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
